# deploy/roles/redis/files/redis-server.conf
description "redis server"

setgid redis
setuid redis

start on runlevel [23]
stop on shutdown
respawn

script
  /usr/bin/redis-server /etc/redis/redis.conf
end script
